# conf.d/novnc_proxy.conf
[program:novnc_proxy]
command=/conf.d/novnc_proxy.sh
autorestart=true
